# supervisord.conf
[program:node]
environment = NODE_ENV=production
command = node app.js
directory = /home/dotcloud/current

